# infra/podman/quadlet/portfolio.container
[Unit]
Description=Portfolio Django app
After=postgres.service portfolio_media-volume.service
Requires=postgres.service portfolio_media-volume.service

[Container]
ContainerName=portfolio
Image=ghcr.io/arsham909/my_portfolio:latest
Network=portfolio.network
EnvironmentFile=/etc/portfolio/portfolio.env
Environment=SQL_HOST=postgres
Environment=SQL_PORT=5432

# User uploads (admin ImageField/FileField) — named volume survives restarts.
Volume=portfolio_media.volume:/home/app/web/mediafiles:Z

PublishPort=127.0.0.1:8000:8000

HealthCmd=python -c "import urllib.request,sys; sys.exit(0 if urllib.request.urlopen('http://localhost:8000/healthz/',timeout=2).status==200 else 1)"
HealthInterval=30s
HealthRetries=3
HealthTimeout=5s
HealthStartPeriod=60s
NoNewPrivileges=true

[Service]
Restart=always
TimeoutStartSec=240

[Install]
WantedBy=default.target
