# infra/podman/quadlet/traefik.container
[Unit]
Description=Traefik v3 reverse proxy
After=network-online.target letsencrypt-volume.service
Wants=network-online.target
Requires=letsencrypt-volume.service

[Container]
ContainerName=traefik
Image=docker.io/library/traefik:v3.1
Network=host
Volume=/etc/traefik/traefik.yml:/etc/traefik/traefik.yml:ro,Z
Volume=/etc/traefik/dynamic:/etc/traefik/dynamic:ro,Z
Volume=letsencrypt.volume:/letsencrypt:Z
EnvironmentFile=/etc/portfolio/traefik.env
ReadOnly=true
Tmpfs=/tmp
Tmpfs=/run
NoNewPrivileges=true
DropCapability=ALL
AddCapability=NET_BIND_SERVICE

[Service]
Restart=always
TimeoutStartSec=60
MemoryHigh=128M
MemoryMax=192M

[Install]
WantedBy=multi-user.target
