# infra/backup/restic-backup.timer
[Unit]
Description=Weekly restic backup timer

[Timer]
# Sunday 03:30 UTC.
OnCalendar=Sun *-*-* 03:30:00
# If box was off at trigger, run on next boot.
Persistent=true
# Jitter 0-10min to avoid thundering herd if more boxes added later.
RandomizedDelaySec=10min

[Install]
WantedBy=timers.target
